# aria2/@svc/aria2.service
[Unit]
Description=Aria2 Daemon

[Service]
ExecStart=/usr/bin/aria2c
OOMScoreAdjust=-1000

[Install]
WantedBy=default.target
