# systemd/udiskie.service
[Unit]
Description=Automounter for removable media

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
