# systemd/udiskie.service
[Unit]
Description=Automounter for removable media

[Service]
ExecStart=/usr/bin/udiskie
OOMScoreAdjust=-1000

[Install]
WantedBy=default.target
